# wrapper-nvidia.service
[Unit]
Description=wrapper-nvidia Node.js proxy (NVIDIA NIM)
After=network.target
Documentation=https://github.com/lokah1945/wrappers
# Limit rapid restart attempts (was in [Service] before — wrong section).
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=/root/wrapper/nvidia
ExecStart=/usr/bin/node /root/wrapper/nvidia/src/index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/root/wrapper/nvidia/.env
Environment=LISTEN_PORT=9100
# Also bind 9910 on the SAME http.Server so Claude Code's gateway model
# discovery URL (CLAUDE_CODE_GATEWAY_MODEL_DISCOVERY_URL=http://localhost:9910/v1/models)
# keeps working without editing client configs. Hermes ILMA + Claude base URL use 9100.
Environment=ANY_ALSO_PORTS=9910
Environment=METRICS_DB=/root/wrapper/nvidia/metrics.db
# Tight keepalive/headers timeouts detect dead sockets promptly, but the OVERALL
# server socket lifetime (serverInstance.timeout, derived from
# ANTI_SILENCE_TIMEOUT_MS in src/index.js) MUST outlast the longest legitimate
# silent 'think' of a large reasoning model (nemotron-3-ultra, deepseek-v4-pro,
# glm-5.x hold the socket 60-150s+ before the first token). The code enforces
#   serverInstance.timeout = max(ANTI_SILENCE_TIMEOUT_MS|SERVER_REQUEST_TIMEOUT_MS,
#                                TTFT_TIMEOUT_MS+30000, 60000)
# and STREAM_REQUEST_TIMEOUT_SEC defaults to 900s, so the invariant is
#   ANTI_SILENCE_TIMEOUT_MS >= 900s + 60s = 960000ms.
# The previous ANTI_SILENCE_TIMEOUT_MS=45000 override VIOLATED this invariant and
# silently 504'd healthy large-reasoning-model streams, the exact premature-kill
# risk called out for model (b). We therefore DO NOT override it: the documented
# 960000ms default in .env.example (>= STREAM_REQUEST_TIMEOUT_SEC + 60s) applies.
# SERVER_KEEPALIVE_TIMEOUT_MS=65000 provides reasonable connection reuse while
# allowing long reasoning sessions to complete.
Environment=SERVER_REQUEST_TIMEOUT_MS=60000
Environment=SERVER_KEEPALIVE_TIMEOUT_MS=65000
Environment=SERVER_HEADERS_TIMEOUT_MS=35000
TimeoutStopSec=15
# Resource limits — cap memory so a runaway upstream response can't OOM
# the host. Conservative; tune up if needed.
MemoryMax=512M
TasksMax=512

[Install]
WantedBy=multi-user.target
