# opencode/wrapper-opencode.service
[Unit]
Description=wrapper-opencode: OpenCode API proxy (port 9103)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/wrapper/opencode
Environment=PYTHONPATH=/root/wrapper/opencode
EnvironmentFile=-/root/wrapper/opencode/.env
Environment=LOG_FILE=/root/wrapper/opencode/opencode.log
ExecStart=/usr/bin/python3 -m uvicorn src.main:app --host 127.0.0.1 --port 9103
Restart=always
RestartSec=3
StandardOutput=append:/root/wrapper/opencode/opencode.log
StandardError=append:/root/wrapper/opencode/opencode.log

[Install]
WantedBy=default.target
