# opencode/wrapper-opencode.service
[Unit]
Description=wrapper-opencode - OpenCode API Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/wrapper/opencode
Environment=PYTHONPATH=/root/wrapper/opencode
EnvironmentFile=/root/wrapper/opencode/.env
ExecStart=/usr/local/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 9107
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
