# nous/systemd/wrapper-nous.service
[Unit]
Description=wrapper-nous: Anthropic+Responses+Chat proxy for Nous (port 9107)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/wrapper/nous/wrapper_nous.py
Restart=always
RestartSec=3
WorkingDirectory=/root/wrapper/nous
EnvironmentFile=-/root/wrapper/nous/.env
StandardOutput=append:/root/wrapper/nous/wrapper_nous.log
StandardError=append:/root/wrapper/nous/wrapper_nous.log

[Install]
WantedBy=default.target
