# nous/wrapper-nous.service
[Unit]
Description=wrapper-nous v2.0.0 (FastAPI) — 100/100 OpenAI + Anthropic compatible proxy for Nous Research
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/root/wrapper/nous/.env
ExecStart=/usr/bin/python3 -m uvicorn wrapper_nous:app --host 127.0.0.1 --port 9106 --workers 1
Restart=always
RestartSec=3
WorkingDirectory=/root/wrapper/nous
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
