# nvidia-python/systemd/wrapper-nvidia-python.service
[Unit]
Description=wrapper-nvidia-python: NVIDIA NIM FastAPI proxy (port 9101)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/wrapper/nvidia-python
Environment=PYTHONPATH=/root/wrapper/nvidia-python
EnvironmentFile=-/root/wrapper/nvidia-python/.env
Environment=LOG_FILE=/root/wrapper/nvidia-python/nvidia_py.log
ExecStart=/usr/bin/python3 -m uvicorn src.main:app --host 127.0.0.1 --port 9101
Restart=always
RestartSec=3
StandardOutput=append:/root/wrapper/nvidia-python/nvidia_py.log
StandardError=append:/root/wrapper/nvidia-python/nvidia_py.log

[Install]
WantedBy=default.target
