# nvidia_backup_20260629_131053/wrapper-nvidia.service
[Unit]
Description=wrapper-nvidia Node.js proxy (NVIDIA NIM)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/wrapper/nvidia
ExecStart=/usr/bin/node /root/wrapper/nvidia/src/index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/root/wrapper/nvidia/.env
Environment=LISTEN_PORT=9100
Environment=METRICS_DB=/root/wrapper/nvidia/metrics.db
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
